# board/ieiworld/qc-additions/service/qc_test.service
[Unit]
Description=qc test

[Service]
ExecStart=/qc/main.sh

[Install]
WantedBy=default.target
