# board/ieiworld/test-additions/service/reboot_test.service
[Unit]
Description=reboot test

[Service]
ExecStart=/usr/bin/reboot_test.sh

[Install]
WantedBy=default.target
